# cdb.cron
# Hard code the ruby file path on production for now
# MIN HOUR DAY MONTH DAY-OF-WEEK COMMAND
0 3 * * * ruby -rubygems ~/MyCardRoid-API/app/cdb.rb >  ~/MyCardRoid-API/log/cdb.log 2>&1
