# cdb.cron
# Hard code the ruby file path on production for now
# MIN HOUR DAY MONTH DAY-OF-WEEK COMMAND
0 3 * * * cd ~/msk86/MyCardRoid-API/ && ruby -rubygems app/cdb.rb > log/cdb.log 2>&1
